# shell-script/crontab/schedule-todo-minuto
* * * * * echo "test" >> /home/cron-exemplo.txt
